# data/bombyx.service.in
[Unit]
Description=bombyx network manager
Wants=network.target
After=network-pre.target dbus.service
Before=network.target @DISTRO_NETWORK_SERVICE@

[Service]
Type=dbus
BusName=org.fpemud.Bombyx
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=@sbindir@/bombyx --no-daemon
Restart=on-failure
CapabilityBoundingSet=CAP_NET_ADMIN CAP_DAC_OVERRIDE CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SETGID CAP_SETUID CAP_SYS_MODULE CAP_AUDIT_WRITE CAP_KILL CAP_SYS_CHROOT
ProtectSystem=true
ProtectHome=read-only

[Install]
WantedBy=multi-user.target

# We want to enable bombyx-wait-online.service whenever this service
# is enabled. bombyx-wait-online.service has
# WantedBy=network-online.target, so enabling it only has an effect if
# network-online.target itself is enabled or pulled in by some other unit.
Also=bombyx-wait-online.service
